# custom/cfg/global/crontab/zpool
30 10 * * * zpool status -x | grep -v 'healthy'
